# client-daemon/service/linux/goconnect-daemon.service
[Unit]
Description=GoConnect Daemon
[Service]
ExecStart=/usr/local/bin/goconnect-daemon
Restart=on-failure
[Install]
WantedBy=multi-user.target
